# etc/rc.local
#!/bin/sh -e

# on the first container start we need to create SSL keys
# TODO update if exists, but expired
CERT1=/etc/ninjam/privkey.pem
CERT2=/etc/ninjam/fullchain.pem
if [[ ! -f "$CERT1" || ! -f "$CERT2" ]]; then
    /etc/ninjam/letsencrypt.sh
fi

service php7.0-fpm start
service supervisor start
sleep 5
service icecast2 start

/usr/sbin/nginx -c /etc/nginx/nginx.conf
